# ubuntu.deprecated/rstudio/supervisord.conf
[supervisord]
nodaemon=true

[program:userconf]
command=/usr/bin/userconf.sh

[program:rserver]
command=/usr/lib/rstudio-server/bin/rserver
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log

